# ucsschool-umc-printermoderation/cron.daily/printermoderation_cleanup
#!/bin/sh
CACHEDIR=/var/cache/printermoderation

find ${CACHEDIR} -P | sed -e '1d' | xargs rm -rf
