# sfb-freeradius/debian/postinst
#! /bin/sh
#
# Univention Squid
#  postinst script
#
# Copyright (C) 2004, 2005, 2006 Univention GmbH
#
# http://www.univention.de/
#
# All rights reserved.
#
# This program is free software; you can redistribute it and/or modify
# it under the terms of the GNU General Public License version 2 as
# published by the Free Software Foundation.
#
# Binary versions of this file provided by Univention to you as
# well as other copyrighted, protected or trademarked materials like
# Logos, graphics, fonts, specific documentations and configurations,
# cryptographic keys etc. are subject to a license agreement between
# you and Univention.
#
# This program is distributed in the hope that it will be useful,
# but WITHOUT ANY WARRANTY; without even the implied warranty of
# MERCHANTABILITY or FITNESS FOR A PARTICULAR PURPOSE.  See the
# GNU General Public License for more details.
#
# You should have received a copy of the GNU General Public License
# along with this program; if not, write to the Free Software
# Foundation, Inc., 51 Franklin St, Fifth Floor, Boston, MA  02110-1301  USA

set -e

#DEBHELPER#

case "$1" in
    configure)

		eval `univention-config-registry shell`
		univention-config-registry set \
			freeradius/conf/log/strippednames?yes \
			freeradius/conf/log/auth?yes \
			freeradius/conf/log/auth/badpass?yes \
			freeradius/conf/log/auth/goodpass?no \
			freeradius/firewall/server/start?2 \
			freeradius/firewall/server/end?19 \
			freeradius/firewall/privateip/start?240 \
			freeradius/firewall/privateip/end?254 \
			freeradius/firewall/enable/internalserver?false \
			freeradius/firewall/enable/broadcasts?false \
			freeradius/firewall/enable/radiusserver?false \
			freeradius/firewall/dnsserver?10.108.1.10 \
			freeradius/bridge/interface/name?br0 \
			freeradius/bridge/interface/primary?eth0 \
			freeradius/bridge/interface/secondary?eth1 \


		/usr/lib/sfb-freeradius/setupCA.sh

		update-rc.d sfb-wlan-bridge defaults 12 75 || true

		/etc/init.d/univention-iptables restart

		/etc/init.d/sfb-wlan-bridge restart || true

		/etc/init.d/freeradius restart
		;;

    abort-upgrade|abort-remove|abort-deconfigure)
		;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 0
		;;
esac


exit 0
